# acmed.service
[Unit]
Description=The ACME Daemon

[Service]
ExecStart=/opt/acmed/index.js
Restart=always

[Install]
WantedBy=multi-user.target
